# barch.service
[Unit]
Description=sdbarker.com Chiliproject

[Service]
User=root
ExecStart=/etc/profile.d/custom.sh

[Install]
WantedBy=multi-user.target
